# Client/startup.service
[Unit]
Description=Setup Startup
[Service]
ExecStart=/etc/startup/startup
Restart=no
[Install]
WantedBy=multi-user.target