# MySQLDatabase/startup.service
[Unit]
Description=Setup Startup
[Service]
ExecStart=/etc/startup/startup
[Install]
WantedBy=multi-user.target